# scripts/agentdesk-dcserver.service
[Unit]
Description=AgentDesk Discord Control Server
After=network.target

[Service]
Type=simple
Environment=AGENTDESK_ROOT_DIR=%h/.adk/release
ExecStart=%h/.adk/release/bin/agentdesk dcserver
WorkingDirectory=%h/.adk/release
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.adk/release/logs/dcserver.stdout.log
StandardError=append:%h/.adk/release/logs/dcserver.stderr.log

[Install]
WantedBy=default.target
